# rskj-ubuntu-installer/rskj_package_trusty/debian/rskj.postrm
#!/bin/sh -e
# Source debconf library.
. /usr/share/debconf/confmodule
 
# Exit on error
set -e

update-rc.d -f rsk remove

db_purge

#DEBHELPER#
